# image/firstboot/neural-ice-seed-import.service
[Unit]
Description=Neural ICE — import the preloaded image seed into podman storage (first boot)
Documentation=https://github.com/Neural-ICE/ICE-CoreOS
# Needs the encrypted data volume (where the installer staged the seed) mounted first.
RequiresMountsFor=/var/lib/neural-ice/data
# Only when a seed is actually present (PRELOADED edition); instant no-op otherwise (LIGHT).
ConditionPathExists=/var/lib/neural-ice/data/seed/images
# Run before the appliance workloads that would otherwise pull images. App Quadlets
# (ICE-AC1) should also carry `After=neural-ice-seed-import.service` for a hard ordering.
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/neural-ice-seed-import.sh
# Large local imports (vLLM ~19 GB) — generous timeout.
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
