# image/tui/neural-ice-tui.service
[Unit]
Description=Neural ICE — appliance console dashboard (TUI on tty1)
Documentation=https://github.com/Neural-ICE/ICE-CoreOS
After=network-online.target
Wants=network-online.target
# The TUI owns tty1; it must not fight a login getty for the console.
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
# Silence kernel console printing on tty1 so dmesg lines don't fight the TUI (debug builds
# set console=tty1). Re-enabled on stop. '+' runs with full privileges (past the sandbox).
ExecStartPre=+/usr/bin/dmesg -D
ExecStartPre=/bin/sh -c 'exec </dev/tty1 >/dev/tty1 2>&1; TERM=linux setterm --clear all 2>/dev/null; clear'
ExecStart=/usr/local/bin/neural-ice-tui
ExecStopPost=+/usr/bin/dmesg -E
Restart=always
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5

# Bind to the physical console tty1.
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty1
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no

# Hardening: the dashboard only READS state (proc, state files, nvidia-smi, systemctl,
# chronyc). No writes to the OS. NoNewPrivileges keeps it from escalating.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
